# systemd/tg-ws-proxy.service
[Unit]
Description=Telegram WebSocket bridge proxy
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
WorkingDirectory=%h/Telegram/tg_proxy
ExecStart=%h/Telegram/tg_proxy/run_proxy.sh run
Restart=on-failure
RestartSec=2

[Install]
WantedBy=default.target
