# systemd/tg-ws-proxy.service
[Unit]
Description=Telegram WebSocket bridge proxy
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
WorkingDirectory=%h/Telegram/tg_proxy
ExecStart=/home/alex/miniconda3/envs/p313/bin/python %h/Telegram/tg_proxy/tg_ws_proxy.py run
Restart=on-failure
RestartSec=2

[Install]
WantedBy=default.target
